# linux/etc/init/osrm.conf
# osrm

description "Open Source Routing Machine"
author "delawen <marias@emergya.com>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Start the process

#script
#  echo "Starting OSRM..."
#  cd /home/devel/Project-OSRM/
#  ./osrm-routed &
#  echo "OSRM started"
#  cd $pwd
#end script

exec osrm

